# weather.cron
*/10 * * * * root /home/pi/repos/weather/weather.sh /var/www/html/images/
